# airootfs/etc/profile.d/env-path.sh
# /sbin is not used on ArchLinux but it is often required in chroot
export PATH=${PATH}:/sbin:/usr/share/sysrescue/bin/
