# airootfs/etc/systemd/system/getty@.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --login-options '-p -f -- \\u' --noclear %I $TERM
