# airootfs/etc/systemd/system/sysrescue-initialize.service
[Unit]
Description=SystemRescue Initialization
Before=getty-pre.target
Wants=getty-pre.target

[Service]
Type=oneshot
ExecStart=/etc/systemd/scripts/sysrescue-initialize.py
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
