# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --keep-baud -l /usr/bin/bashlogin --noclear %I linux
