# airootfs/etc/systemd/system/serial-getty@.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --login-options '-p -f -- \\u' --keep-baud 115200,57600,38400,9600 %I $TERM
